# kc868/kc868_iluminacao_status.conf
[program:kc868_iluminacao_status.php]
command=php kc868_iluminacao_status.php
directory=/var/www/html/Alarme_MQTT/kc868
user=epaminondas
stdout_logfile=/var/www/html/Alarme_MQTT/kc868/kc868_iluminacao.log
redirect_stderr=true
autostart=true
autorestart=true
startsecs=0